# elastix/drbd/sync/sync_slave.cron
*/5 * * * *  root /usr/local/sbin/sync_slave > /dev/null 2>&1